# demo-leader/scripts/systemd/demo-leader-api.service
[Unit]
Description=ACPs demo-leader API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Default deployment root. Adjust both paths if you install elsewhere.
WorkingDirectory=/opt/demo-leader
ExecStart=/bin/bash /opt/demo-leader/scripts/start-leader-api.sh
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=demo-leader-api
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Uncomment after creating a dedicated service account.
# User=demo-leader
# Group=demo-leader

# The application loads .env from WorkingDirectory via python-dotenv,
# so this unit intentionally does not use EnvironmentFile=.

[Install]
WantedBy=multi-user.target