# demo-partner/scripts/systemd/demo-partner.service
[Unit]
Description=ACPs demo-partner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Default deployment root. Adjust both paths if you install elsewhere.
WorkingDirectory=/opt/demo-partner
ExecStart=/opt/demo-partner/.venv/bin/python -m partners.main
Environment=PYTHONUNBUFFERED=1
Environment=PARTNERS_ONLINE_DIR=/opt/demo-partner/partners/online
SyslogIdentifier=demo-partner
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Uncomment after creating a dedicated service account.
# User=demo-partner
# Group=demo-partner

# The application loads .env from WorkingDirectory via python-dotenv,
# so this unit intentionally does not use EnvironmentFile=.

[Install]
WantedBy=multi-user.target