# registry-server/scripts/systemd/registry-server.service
[Unit]
Description=ACPs Registry Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Default deployment root. Adjust both paths if you install elsewhere.
WorkingDirectory=/opt/registry-server
ExecStart=/opt/registry-server/.venv/bin/python -m app.runtime_dual_listener
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=registry-server
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Uncomment after creating a dedicated service account.
# User=registry-server
# Group=registry-server

# The application loads .env from WorkingDirectory via pydantic-settings,
# so this unit intentionally does not use EnvironmentFile=.

[Install]
WantedBy=multi-user.target