# discovery-server/scripts/systemd/discovery-server.service
[Unit]
Description=ACPs Discovery Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Default deployment root. Adjust both paths if you install elsewhere.
WorkingDirectory=/opt/discovery-server
ExecStart=/opt/discovery-server/.venv/bin/python -m app.main
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=discovery-server
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

# Uncomment after creating a dedicated service account.
# User=discovery-server
# Group=discovery-server

# The application loads .env from WorkingDirectory via pydantic-settings,
# so this unit intentionally does not use EnvironmentFile=.

[Install]
WantedBy=multi-user.target
